# modules/zookeeper/templates/zookeeper.service.erb
# Autogenerated by Puppet from zookeeper/templates/zookeeper.service.erb
# Do not edit

[Unit]
Description=Apache ZooKeeper
<%if ! [nil, :undefined, :undef].include?(scope.lookupvar("zookeeper::systemd_unit_want")) -%>
Wants=<%= scope.lookupvar("zookeeper::systemd_unit_want")%>
<% end -%>
<%if ! [nil, :undefined, :undef].include?(scope.lookupvar("zookeeper::systemd_unit_after")) -%>
After=<%= scope.lookupvar("zookeeper::systemd_unit_after") %>
<% end -%>

[Service]
ExecStart=/bin/bash -c "source /etc/zookeeper/conf/environment; CLASSPATH=/usr/share/zookeeper/zookeeper.jar:/usr/share/zookeeper/lib/*:$CLASSPATH; CLASSPATH=$([ -r /usr/share/zookeeper/bin/zkEnv.sh ] && . /usr/share/zookeeper/bin/zkEnv.sh ; echo $CLASSPATH); mkdir -p /var/log/zookeeper; /usr/bin/java -Dzookeeper.log.dir=/var/log/zookeeper -Dzookeeper.root.logger=$ZOO_LOG4J_PROP -cp /usr/share/zookeeper/zookeeper.jar:/usr/share/zookeeper/lib/*: org.apache.zookeeper.server.quorum.QuorumPeerMain /etc/zookeeper/conf/zoo.cfg"
Restart=always
RemainAfterExit=no
SyslogIdentifier=zookeeper
User=<%= scope.lookupvar("zookeeper::user") %>
Group=<%= scope.lookupvar("zookeeper::group") %>

[Install]
WantedBy=multi-user.target
