# modules/swift/files/devicecheck.cron
* * * * * root /usr/local/bin/disable_rsync.py
